# profile.d/z-000-init.csh
if ( ! $?CC_CLUSTER ) then
	setenv CC_CLUSTER "computecanada"
endif
