# profile.d/z-30-easybuild.sh
#source /cvmfs/soft.computecanada.ca/easybuild/framework/minimal_bash_completion.bash
#source /cvmfs/soft.computecanada.ca/easybuild/framework/optcomplete.bash
#source /cvmfs/soft.computecanada.ca/easybuild/framework/eb_bash_completion.bash
